# etc/supervisor/conf.d/cupsd.conf
[program:cupsd]
priority=10
command=/usr/sbin/cupsd -c /etc/cups/cupsd.conf -f
user=root
autostart=true
autorestart=true
startsecs=10
stopsignal=QUIT
redirect_stderr=true
